# deploy/systemd/nginx.service.d/override.conf
# nginx self-heal drop-in for the ECS Host (Ubuntu 24.04, systemd 255).
#
# Stock Ubuntu ships nginx.service with NO Restart= policy, so a single failed
# start leaves it dead until a human intervenes. On 2026-06-12 a transient DNS
# failure during the daily unattended-upgrade made `nginx -t` (ExecStartPre) exit
# non-zero, the unit failed before the forking master ever started, and nginx
# stayed DOWN FOR 3 DAYS.
#
# Two independent fixes guard against a recurrence; keep BOTH:
#   1. PRIMARY (load-bearing): deploy/nginx/sentry-tunnel.conf resolves its Sentry
#      upstream LAZILY (resolver + $variable in proxy_pass) so `nginx -t` no longer
#      fails on a DNS blip — the root cause of the 2026-06-12 ExecStartPre failure.
#      With it in place, the startup-fatal DNS mode cannot recur at all.
#   2. SAFETY NET (this file): if ExecStartPre or the master process ever fails for
#      ANY other reason, systemd keeps restarting nginx until the cause clears.
#
# Type=forking + ExecStartPre interaction: when ExecStartPre (/usr/sbin/nginx -t)
# exits non-zero, systemd marks the unit FAILED *before* forking the master, which
# Restart=on-failure DOES cover (on-failure = nonzero exit of any Exec*, the main
# process, signals, and watchdog).
#
# Why StartLimitIntervalSec=0 (retry forever) and NOT a bounded burst: a bounded
# burst (e.g. 30 tries x 10s ~= 5 min) LATCHES the unit into a permanent
# 'failed (start-limit-hit)' state once exhausted, needing `systemctl reset-failed`
# or a reboot — it would NOT self-heal an outage longer than ~5 min (exactly the
# multi-hour / 3-day scenario). IntervalSec=0 disables the start-rate limiter, so
# Restart=on-failure retries indefinitely at a 10s cadence and self-heals the moment
# the cause clears, at negligible cost on an idle box (journald rate-limits the log
# noise). A deliberate `systemctl stop` exits 0 and is still honored.
#
# Verify after install:
#   systemctl show nginx -p Restart -p RestartUSec -p StartLimitIntervalUSec
#   # expect: Restart=on-failure  RestartUSec=10s  StartLimitIntervalUSec=0
[Service]
Restart=on-failure
RestartSec=10s

[Unit]
# 0 disables the start-rate limiter -> retry forever (see header). Do NOT set this
# to a bounded window; a bounded burst parks the unit dead after ~5 min.
StartLimitIntervalSec=0
